# rootfs/etc/services.d/sphinx/run
#!/usr/bin/with-contenv bash

echo "Starting searchd"

exec searchd --console --pidfile --config /etc/sphinx/sphinx.conf
